# scripts/com.homelab.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.homelab.backup</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/opt/homelab/scripts/backup.sh</string>
  </array>

  <!-- Inject PATH so launchd can find rclone and docker -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/yourusername</string>
  </dict>

  <!-- Run every day at 02:00 -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>2</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/yourusername/Library/Logs/homelab-backup.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yourusername/Library/Logs/homelab-backup-error.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
